# server/deploy/systemd/fenix-worker.service
[Unit]
Description=Fenix background worker
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fenix/server
Environment=NODE_ENV=production
Environment=PUBLIC_BASE_URL=https://garphoenix.com
ExecStart=/usr/bin/env npm run worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
